# backend/symbioflows.service
[Unit]
Description=Revolutionary AI Matching System
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/your/backend
Environment=PATH=/path/to/your/venv/bin
ExecStart=/path/to/your/venv/bin/python production_startup.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
